# services/yourss@.timer
[Unit]
Description=Timer for running yourss@%i.service

[Timer]
OnCalendar=daily

[Install]
WantedBy=timers.target
